# configs/rc.local
#!/bin/sh -e

/root/setup_networking.sh &
/root/health_logs.sh &
/root/create_templates.sh &
sleep 2m
/root/recycle.sh 101 > /root/log.txt 2>&1 &
