# configs/rc.local
#!/bin/sh -e

/root/setup_networking.sh &
/root/create_templates.sh &
sleep 1.5m
/root/recycle.sh 101 201 > /root/log1.txt 2>&1 &
/root/recycle.sh 102 201 > /root/log2.txt 2>&1 &
/root/recycle.sh 103 201 > /root/log3.txt 2>&1 &
